# deploy/pi/geo-dashboard.service
[Unit]
Description=geoanalytics dashboard + REST API (FastAPI :8800) — Raspberry Pi
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/News
Environment=PYTHONPATH=%h/News
Environment=GEO_DB_HOST=127.0.0.1
# Ollama (Qwen для ask/нарративов) живёт на ноуте; когда ноут выключен — LLM-фичи дашборда
# деградируют (read-only/числовые страницы работают; advisory-lock/фолбэк отрабатывает недоступность).
Environment=GEO_OLLAMA_HOST=http://192.168.0.13:11434
# 0.0.0.0 → дашборд доступен по LAN (http://192.168.0.114:8800/) даже когда ноут выключен.
ExecStart=%h/News/.venv/bin/geo serve --host 0.0.0.0 --port 8800
MemoryAccounting=yes
MemoryHigh=1500M
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
